# Backup_configs/old_crawlers/tempCrawler.conf
# Run sysCrawler.py script

start on runlevel [2345]
stop on runlevel [!2345]

respawn
#exec bash /home/pi/DataCrawler/tempCrawler.sh




