# Backup_configs/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

# Print the IP address
#_IP=$(hostname -I) || true
#if [ "$_IP" ]; then
(sleep 10; python /home/pi/Scripts/Whats_My_IP/whats_my_ip.py)&
#fi
(sleep 12; python /home/pi/Scripts/Gesture_Control/UDP_Sensor_Rx.py)& > /dev/null
(sleep 13; python /home/pi/Scripts/Smart_DoorBell/WaitDoorbell.py)& >/dev/null
(sleep 15; python /home/pi/Scripts/Sensor_Loggers/Indoor_Temp_Logger.py)& >/dev/null
(sleep 16; python /home/pi/Scripts/TV_Proximity_Sensor/proximity_sensor_v2.py)& >/dev/null
exit 0
